# encode.conf
start on startup
respawn

setuid dlink
setgid dlink

exec bash /home/dlink/encode.sh

