# record-ip.conf
start on startup
respawn

script

while true
do
wget --spider -q http://viscerallogic.com/ip_address/test.php?query_type=record&name=nt-a3000
sleep 3600
done

end script
